# supervisord.conf
# supervisord configuration for forever-claude-template background services.
#
# `uv run bootstrap` runs first-boot setup and then execs `supervisord -n -c
# supervisord.conf` from the repo root (/mngr/code) in the bootstrap tmux
# window. supervisord then starts and supervises every [program:*] below.
#
# Conventions:
#   - Services inherit the agent environment (MNGR_AGENT_STATE_DIR,
#     CLAUDE_CONFIG_DIR, MNGR_HOST_DIR, MNGR_AGENT_ID, GH_TOKEN, ...) from the
#     bootstrap shell that exec'd supervisord -- there is no per-service
#     `environment=` enumeration.
#   - All commands run from the repo root (directory=/mngr/code).
#   - Commands that chain with `&&` (or set inline env) are wrapped in
#     `bash -c "..."` because supervisord exec's the command directly (no shell).
#   - Logs are container-local under /var/log/supervisor (created by bootstrap),
#     rotated by supervisord, with separate stdout/stderr files per service.
#     They are deliberately NOT under runtime/ -- they are not backed up.
#
# To add/remove/change a service, edit this file and run
# `supervisorctl reread && supervisorctl update` (see the edit-services skill).

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
# We intentionally run as root (the agent container/VM has no separate service
# user); set user=root so supervisord doesn't emit a CRIT about not dropping privs.
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

# Default knobs shared by every [program:*] section below. supervisord has no
# global program defaults, so each program restates these.
#   autostart=true              start on supervisord boot
#   autorestart=true            restart a long-lived daemon whenever it exits
#   startretries=...            cap crash-loop retries per supervisord run
#   stopasgroup/killasgroup     signal the whole bash -c process group on stop
#   stdout/stderr_logfile       separate, rotated, container-local log files

[program:system_interface]
command=bash -c "python3 scripts/forward_port.py --url http://localhost:8000 --name system_interface && system-interface"
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/system_interface-stdout.log
stderr_logfile=/var/log/supervisor/system_interface-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

[program:web]
command=bash -c "python3 scripts/forward_port.py --url http://localhost:8080 --name web && uv run web-server"
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/web-stdout.log
stderr_logfile=/var/log/supervisor/web-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

[program:cloudflared]
command=uv run cloudflare-tunnel
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/cloudflared-stdout.log
stderr_logfile=/var/log/supervisor/cloudflared-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

[program:app-watcher]
command=uv run app-watcher
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/app-watcher-stdout.log
stderr_logfile=/var/log/supervisor/app-watcher-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# The cron daemon, in the foreground so supervisord owns it. Runs the
# /etc/cron.d/ entries: precise / sub-daily schedule lines, plus the
# every-minute daily-job checker lines (scripts/run_daily_job.sh) that decide
# when daily jobs like the Caretaker actually fire (at the due hour when the
# container is up, or the first minute it is back up after a missed day).
[program:cron]
command=/usr/sbin/cron -f
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/cron-stdout.log
stderr_logfile=/var/log/supervisor/cron-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# Periodically commits + pushes runtime/ to mindsbackup/$MNGR_AGENT_ID so that
# Claude memory, ticket state, transcripts, etc. survive container loss. The
# runtime/ worktree is created by bootstrap before supervisord starts, so this
# service always finds it in place. See libs/runtime_backup/README.md.
[program:runtime-backup]
command=uv run runtime-backup
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/runtime-backup-stdout.log
stderr_logfile=/var/log/supervisor/runtime-backup-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# Continuous restic-based backup of the full host_dir (/mngr/) to a remote
# repository. runtime/secrets/restic.env is injected by the minds app;
# runtime/backup.toml is optional user settings. This [program] block, the
# root pyproject registration, and the `uv run host-backup` entry point are a
# stable contract relied on by minds backup-service updates -- do not rename.
# See libs/host_backup/README.md.
[program:host-backup]
command=uv run host-backup
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/host-backup-stdout.log
stderr_logfile=/var/log/supervisor/host-backup-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# Terminal-over-web (ttyd). Registers its port via forward_port.py, then execs
# ttyd. Previously an extra_window; now a supervised service like the rest.
[program:terminal]
command=bash scripts/run_ttyd.sh
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/terminal-stdout.log
stderr_logfile=/var/log/supervisor/terminal-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# One-shot, idempotent install of deferred packages (currently Playwright +
# Chromium + apt libs) on first container boot. autorestart=false + exitcodes=0
# make it run exactly once per boot and stay stopped after a clean exit; it
# skips its own work on subsequent boots via per-package marker files under
# /var/lib/minds/deferred-install/. See scripts/deferred_install.sh.
[program:deferred-install]
command=bash scripts/deferred_install.sh
directory=/mngr/code
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/deferred-install-stdout.log
stderr_logfile=/var/log/supervisor/deferred-install-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# Live-browser web service: a headless Chromium streamed via CDP screencast with
# CDP input + a browser-use agent, at /service/browser/. No virtual display needed.
[program:browser]
command=bash -c "ROOT_PATH=/service/browser python3 scripts/forward_port.py --url http://localhost:8081 --name browser && uv run browser-service"
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/browser-stdout.log
stderr_logfile=/var/log/supervisor/browser-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
