# supervisord.conf
# supervisord configuration for default-workspace-template background services.
#
# `uv run bootstrap` runs first-boot setup and then execs `supervisord -n -c
# supervisord.conf` from the repo root (/mngr/code) in the bootstrap tmux
# window. supervisord then starts and supervises every [program:*] below.
#
# Conventions:
#   - Services inherit the agent environment (MNGR_AGENT_STATE_DIR,
#     CLAUDE_CONFIG_DIR, MNGR_HOST_DIR, MNGR_AGENT_ID, LATCHKEY_*, ...) from
#     the bootstrap shell that exec'd supervisord -- there is no per-service
#     `environment=` enumeration.
#   - All commands run from the repo root (directory=/mngr/code).
#   - Commands that chain with `&&` (or set inline env) are wrapped in
#     `bash -c "..."` because supervisord exec's the command directly (no shell).
#   - Logs are container-local under /var/log/supervisor (created by bootstrap),
#     rotated by supervisord, with separate stdout/stderr files per service.
#     They are deliberately NOT under runtime/ -- they are not backed up.
#
# To add/remove/change a service, edit this file and run
# `supervisorctl reread && supervisorctl update` (see the edit-services skill).

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
# We intentionally run as root (the agent container/VM has no separate service
# user); set user=root so supervisord doesn't emit a CRIT about not dropping privs.
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

# Default knobs shared by every [program:*] section below. supervisord has no
# global program defaults, so each program restates these.
#   autostart=true              start on supervisord boot
#   autorestart=true            restart a long-lived daemon whenever it exits
#   startretries=...            cap crash-loop retries per supervisord run
#   stopasgroup/killasgroup     signal the whole bash -c process group on stop
#   stdout/stderr_logfile       separate, rotated, container-local log files

[program:system_interface]
command=python3 scripts/oom_tag_service.py system_interface bash -c "sh scripts/openhost_stop_placeholder.sh && python3 scripts/forward_port.py --url http://localhost:8000 --name system_interface && system-interface"
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/system_interface-stdout.log
stderr_logfile=/var/log/supervisor/system_interface-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# NOTE: cloudflared, app-watcher, and host-backup are gone on OpenHost: the
# OpenHost router does ingress/auth (no Cloudflare tunnel), there is no desktop
# client consuming service events, and the platform backs up the app-data dir
# (no restic).

# NOTE: GitHub sync ([program:github-sync]) is opt-in and deliberately NOT in
# this default config. The github-sync skill appends its program block (and
# sets up the private repo, git wiring, and runtime/ worktree) when the user
# asks for it. See libs/github_sync/README.md.

# Terminal-over-web (ttyd). Registers its port via forward_port.py, then execs
# ttyd. Previously an extra_window; now a supervised service like the rest.
[program:terminal]
command=python3 scripts/oom_tag_service.py terminal bash scripts/run_ttyd.sh
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/terminal-stdout.log
stderr_logfile=/var/log/supervisor/terminal-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# OOM-prevention daemon. Watches memory and, before the kernel resorts to its
# own arbitrary OOM kill, sheds the worst process (SIGTERM below 10% available
# memory, SIGKILL below 5%; same thresholds for swap). It reads
# /proc/*/oom_score to choose -- which already reflects each process's
# oom_score_adj -- so the startup priority bands (see libs/oom_priority) steer
# its choice, most-expendable first: shared browsers, then agent build/test
# subprocesses, then worker agents, then user agents, then user-created
# services, then the built-in services in strict order (web, app-watcher,
# host-backup, github-sync, cloudflared, system_interface, terminal). --avoid
# shields the never-kill infrastructure that has a distinct process name (sshd,
# supervisord, earlyoom itself, tini, tmux); these keep oom_score_adj 0 and are
# the last resort. cloudflared and ttyd are intentionally NOT in --avoid: they
# carry low service bands so they are shed only after everything above, but are
# not made strictly immortal (both auto-restart). -r 0 silences the periodic
# memory report (kills are still logged). -N records every kill to the shed
# ledger: earlyoom v1.7 execs the -N target directly (not via a shell) and
# disables it at startup unless its first token is an absolute path to an
# *executable* file, so this points straight at the script (which carries a 0755
# mode and a python3 shebang) rather than a "python3 <script>" command line. See
# libs/oom_priority/README.md.
[program:earlyoom]
command=earlyoom -m 10,5 -s 10,5 -r 0 -N "/mngr/code/scripts/earlyoom_record_shed.py" --avoid "^(sshd|supervisord|earlyoom|tini)$|^tmux"
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/earlyoom-stdout.log
stderr_logfile=/var/log/supervisor/earlyoom-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# Backstop for the OOM band tagging (see libs/oom_priority/README.md). On every
# PROCESS_STATE_RUNNING event -- boot and each restart -- it raises the program's
# process (and any children already spawned) up to its expected band, so a
# service whose command omits the oom_tag_service.py prefix lands in the
# user-service band instead of keeping the fully-protected inherited default of
# 0. Raise-only: it never lowers a value, so self-tagged processes (browser) and
# the PROTECTED programs (earlyoom, this listener) are untouched. stdout carries
# the event-listener protocol, so only stderr is logged. buffer_size comfortably
# exceeds the number of programs so a boot-time burst of RUNNING events (queued
# while this listener itself is still starting) can never overflow and drop one.
[eventlistener:oom-tag-backstop]
command=python3 scripts/oom_tag_backstop.py
events=PROCESS_STATE_RUNNING
buffer_size=100
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stderr_logfile=/var/log/supervisor/oom-tag-backstop-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

# One-shot, idempotent install of deferred packages (currently Playwright +
# Chromium + apt libs) on first container boot. autorestart=false + exitcodes=0
# make it run exactly once per boot and stay stopped after a clean exit; it
# skips its own work on subsequent boots via per-package marker files under
# /var/lib/minds/deferred-install/. See scripts/deferred_install.sh.
[program:deferred-install]
command=bash scripts/deferred_install.sh
directory=/mngr/code
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/deferred-install-stdout.log
stderr_logfile=/var/log/supervisor/deferred-install-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

# Live-browser web service: a headless Chromium streamed via CDP screencast with
# CDP input + a browser-use agent, at /service/browser/. No virtual display needed.
#
# OOM: shared browsers are the most expendable thing in the workspace, so this is
# the one service that opts OUT of the protected band. It raises its own
# oom_score_adj to 1000 (the ceiling -- SHARED_BROWSER in libs/oom_priority, one
# above the agent-subprocess band of 900; the oom-tag-backstop listener re-asserts
# it if this inline write failed), and Chromium inherits oom_score_adj across
# fork/exec, so
# every browser the daemon spawns lands at 1000 too. earlyoom picks the highest
# /proc/*/oom_score (= memory badness 0-1000 + oom_score_adj), so a browser always
# outranks even an agent's build/test subprocess and is shed first. Raising one's
# own oom_score_adj is unprivileged; `|| true` keeps a write failure from wedging
# startup (the service just comes up unprotected). The slim daemon parent shares
# the 1000 but holds the least memory, so it's the last of the group earlyoom would
# pick -- and a shed Chromium is dropped from the fleet manifest, so a daemon
# respawn restores nothing.
[program:browser]
command=bash -c "echo 1000 > /proc/self/oom_score_adj || true && ROOT_PATH=/service/browser python3 scripts/forward_port.py --url http://localhost:8081 --name browser && uv run browser-service"
directory=/mngr/code
autostart=true
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/browser-stdout.log
stderr_logfile=/var/log/supervisor/browser-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
